# roles/rclone-backup/templates/rclone-backup.timer.j2
[Unit]
Description=Run rclone backup 4 times a day

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true

[Install]
WantedBy=timers.target
